# build/discovery.service
[Unit]
Description=Bastion - Discovery
After=docker.service
Requires=docker.service nsqd.service

[Service]
User=core
TimeoutStartSec=0
ExecStartPre=-/usr/bin/docker kill %p-%i
ExecStartPre=-/usr/bin/docker rm %p-%i
ExecStartPre=/usr/bin/docker pull quay.io/opsee/bastion:latest
ExecStart=/usr/bin/docker run --name %p-%i --link nsqd:nsqd quay.io/opsee/bastion:latest /discovery
ExecStop=/usr/bin/docker stop %p-%i

[Install]
WantedBy=multi-user.target
