# build/env-lock.path
[Unit]
Description=Bastion - Connector

[Path]
PathExistsGlob=/etc/opsee/bastion-env.sh
