# build/env-lock.service
[Unit]
Description=Triggered after generation of env file
ConditionFileNotEmpty=/etc/opsee/bastion.env

[Service]
Type=oneshot
ExecStart=/bin/sh -c "echo Bastion Env Ready"
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
