# build/env-lock.service
[Unit]
Description=Triggered after generation of env file
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=echo "env file found"
