# build/monitor.service
[Unit]
Description=Bastion - Monitor
After=docker.service
Requires=docker.service nsqd.service

[Service]
User=core
TimeoutStartSec=0
ExecStartPre=-/usr/bin/docker kill %p-%i
ExecStartPre=-/usr/bin/docker rm %p-%i
ExecStartPre=/usr/bin/docker pull quay.io/opsee/bastion:latest
ExecStart=/usr/bin/docker run --name %p-%i -p 4000:4000 --link nsqd:nsqd quay.io/opsee/bastion:latest /monitor
ExecStop=/usr/bin/docker stop %p-%i

[Install]
WantedBy=multi-user.target
