# build/nsqd.service
[Unit]
Description=nsqd
After=docker.service
Requires=docker.service

[Service]
User=core
TimeoutStartSec=0
ExecStartPre=-/usr/bin/docker kill nsqd
ExecStartPre=-/usr/bin/docker rm nsqd
ExecStartPre=/usr/bin/docker pull nsqio/nsq:v0.3.5
ExecStart=/usr/bin/docker run --name nsqd -p 4150:4150 -p 4151:4151 nsqio/nsq /nsqd
ExecStop=/usr/bin/docker stop nsqd

[Install]
WantedBy=multi-user.target
