# build/unhacker.service
[Unit]
Description=Bastion - Security Group Unhacker
Before=shutdown.target reboot.target halt.target

[Service]
Type=oneshot
User=core
TimeoutStopSec=0
EnvironmentFile=/etc/opsee/bastion-env.sh
ExecStart=/bin/true
ExecStop=/usr/bin/docker stop -t 5 hacker
ExecStop=/usr/bin/docker run --name %p quay.io/opsee/bastion:${BASTION_VERSION} /hacker -unhack
ExecStopPost=/usr/bin/docker stop -t 5 %p
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
